# .config/profile.d/desktop.sh
#!/bin/sh

# required by gnome-control-center > 1:3.28
XDG_CURRENT_DESKTOP=GNOME;   export XDG_CURRENT_DESKTOP;


